# robotaxi.service
[Unit]
Description=WaveShare Robotaxi Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/WaveShare
ExecStart=/home/pi/WaveShare/start_robotaxi.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Give the service time to start up
TimeoutStartSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target 